# com.macbook-notify.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macbook-notify.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SERVER_URL</key>
        <string>__SERVER_URL__</string>
        <key>TOPIC</key>
        <string>__TOPIC__</string>
        <key>AUTH_TOKEN</key>
        <string>__AUTH_TOKEN__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/macbook-notify/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/macbook-notify/stderr.log</string>
</dict>
</plist>
